# profile/airootfs/etc/systemd/system/arsenal-selftest.service
[Unit]
Description=Arsenal hardening self-test (boot-gated; CI integration test)
Documentation=https://github.com/BoyUnderThunder/arsenal
# Only runs when the kernel cmdline contains `arsenal.selftest` (set by
# tools/integration-test.sh). On normal boots this condition is false and the
# unit is skipped, so it never affects end users.
ConditionKernelCommandLine=arsenal.selftest
# Run once the hardening is actually applied.
After=apparmor.service nftables.service systemd-sysctl.service systemd-user-sessions.service
Wants=apparmor.service nftables.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/arsenal-selftest
# Mirror output to the serial console so the CI harness can read the sentinel.
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
