# profile/airootfs/etc/systemd/system/serial-getty@ttyS0.service.d/autologin.conf
# Autologin root on the serial console. Used for headless / QEMU boot testing
# (the build's smoke test boots with -nographic and reads ttyS0).
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- \\u' --keep-baud --autologin root 115200,57600,38400,9600 - $TERM
